# delpro-monitor.service
[Unit]
Description=Monitor DelPro database

[Service]
Type=oneshot
User=nobody
Group=nobody
WorkingDirectory=/path/to/
ExecStart=/path/to/update.rb
